# codecheck/codecheck-2366/contrib/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:subot]
command=/usr/bin/python3 /srv/subot/run.py
